# crontab
*/15 * * * * root /opt/render.sh >> /var/log/render.log 2>&1
